# custom/cfg/crontab/zpool-healthy
30 8 * * * zpool status -x | grep -v 'healthy'
